# config/llama-server.service
[Unit]
Description=llama.cpp OpenAI-compatible server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ajmeese7
Group=ajmeese7
WorkingDirectory=/home/ajmeese7
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/lib/wsl/lib
# LD_LIBRARY_PATH needed on WSL2 for CUDA libs; harmless if the path doesn't exist on native Linux.
Environment=LD_LIBRARY_PATH=/usr/lib/wsl/lib
ExecStart=/etc/llama-server/llama-launcher.sh
Restart=always
RestartSec=10
TimeoutStartSec=120
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
