# config/llama-server.service.template
[Unit]
Description=llama.cpp OpenAI-compatible server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=@@USER@@
Group=@@GROUP@@
WorkingDirectory=@@REPO_ROOT@@
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/lib/wsl/lib
# LD_LIBRARY_PATH needed on WSL2 for CUDA libs; harmless on native Linux.
Environment=LD_LIBRARY_PATH=/usr/lib/wsl/lib
ExecStart=@@REPO_ROOT@@/.venv/bin/llms launcher exec --config @@REPO_ROOT@@/config
Restart=always
RestartSec=10
TimeoutStartSec=120
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
