# config/llama-server.v2.service
[Unit]
Description=llama.cpp OpenAI-compatible server (v2 launcher)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ajmeese7
Group=ajmeese7
WorkingDirectory=/home/ajmeese7/src/local-llms
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/lib/wsl/lib
# LD_LIBRARY_PATH needed on WSL2 for CUDA libs; harmless on native Linux.
Environment=LD_LIBRARY_PATH=/usr/lib/wsl/lib
# The v2 launcher is the Python entry point. Phase 7 setup installs it as
# `llms` somewhere on PATH; until then this template documents the intent.
ExecStart=/home/ajmeese7/.local/bin/llms launcher exec --config /home/ajmeese7/src/local-llms/config
Restart=always
RestartSec=10
TimeoutStartSec=120
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
